# 1.0.alpha.404/debian/postinst
#!/bin/sh

set -e

echo "Enabling+Starting Service"

systemctl enable traefik
systemctl start traefik
