# 1.0.alpha.404/debian/prerm
#!/bin/sh

set -e

echo "Disabling Traefik Service"

systemctl stop traefik
systemctl disable traefik
